# config/supervisord.ini
[program:file-share]
directory=/file-share/dist
command=/file-share/dist/file-share /file-share/dist/api.json
user=www

[program:nginx]
command=/usr/sbin/nginx