# config/supervisord.ini
[program:file-share]
directory=/file-share/dist
command=/file-share/dist/file-share /file-share/config/api.json
user=www
stopsignal=INT

[program:nginx]
command=/usr/sbin/nginx

[program:sshd]
command=/usr/sbin/sshd -D
